# packages/merge-steward/infra/merge-steward.service
[Unit]
Description=merge-steward
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=8

[Service]
Type=simple
User=your-user
Group=your-user
WorkingDirectory=/home/your-user

EnvironmentFile=-/home/your-user/.config/merge-steward/runtime.env
EnvironmentFile=-/home/your-user/.config/merge-steward/service.env
Environment=NODE_ENV=production
Environment=PATH=/home/your-user/.local/bin:/usr/local/bin:/usr/bin:/bin

# Service-owned secrets live in encrypted systemd credentials.
LoadCredentialEncrypted=merge-steward-webhook-secret:/etc/credstore.encrypted/merge-steward-webhook-secret.cred
LoadCredentialEncrypted=merge-steward-github-app-pem:/etc/credstore.encrypted/merge-steward-github-app-pem.cred

ExecStart=/usr/bin/env merge-steward serve
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
PrivateMounts=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/your-user/.config/merge-steward /home/your-user/.local/state/merge-steward /home/your-user/.local/share/merge-steward

[Install]
WantedBy=multi-user.target
