# infra/patchrelay.service
[Unit]
Description=PatchRelay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=patchrelay
Group=patchrelay
WorkingDirectory=/opt/patchrelay
EnvironmentFile=/opt/patchrelay/.env
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env npm run start
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/patchrelay /var/lib/patchrelay /srv/patchrelay/worktrees /var/log/patchrelay
StateDirectory=patchrelay
RuntimeDirectory=patchrelay
LogsDirectory=patchrelay

[Install]
WantedBy=multi-user.target
