# infra/patchrelay.service
[Unit]
Description=PatchRelay (systemd user service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/your-user
EnvironmentFile=/home/your-user/.config/patchrelay/.env
Environment=NODE_ENV=production
Environment=PATCHRELAY_CONFIG=/home/your-user/.config/patchrelay/patchrelay.yaml
Environment=PATH=/home/your-user/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env patchrelay serve
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/your-user/.config/patchrelay /home/your-user/.local/state/patchrelay /home/your-user/.local/share/patchrelay

# PatchRelay is intended to run as your real user so Codex inherits your
# existing git, SSH, and local tool permissions.
# Add your managed repository roots to ReadWritePaths if you keep hardening enabled.

[Install]
WantedBy=default.target
