# infra/patchrelay-reload.service
[Unit]
Description=PatchRelay reload helper

[Service]
Type=oneshot
ExecStart=/usr/bin/env systemctl reload-or-restart patchrelay.service
