# infra/patchrelay.path
[Unit]
Description=Watch PatchRelay config and env changes

[Path]
Unit=patchrelay-reload.service
PathChanged=/home/your-user/.config/patchrelay/runtime.env
PathChanged=/home/your-user/.config/patchrelay/service.env
PathChanged=/home/your-user/.config/patchrelay/patchrelay.json
TriggerLimitIntervalSec=5
TriggerLimitBurst=1

[Install]
WantedBy=multi-user.target
